# flatpak-background-update.service
[Unit]
Description=Script to Update Flatpak Packages

[Service]
ExecStart=/usr/bin/flatpak update -y
