# flatpak-update.service
[Unit]
Description=Flatpak Automatic Updates

[Service]
ExecStart=/usr/bin/flatpak update -y

[Install]
WantedBy=multi-user.target
